# files/etc/cron.d/acd_cli_sync
* * * * * root /usr/local/sbin/acd_cli_sync

