# files/etc/supervisor/services-available/mount_media.conf
[program:mount_media]
command=/usr/local/sbin/mount_media
priority=10
autostart=true
autorestart=false
startretries=0

