# files/etc/supervisor/services-available/plexmediaserver.conf
[program:plexmediaserver]
command=/usr/local/sbin/plexmediaserver
priority=500
user=plex
autorestart=true
autostart=true
stopasgroup=true
killasgroup=true

